# deploy/larynx-cleanup.service
[Unit]
Description=Larynx TTS Audio Cleanup Service
After=network.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/larynx/backend
EnvironmentFile=/var/www/larynx/backend/.env
ExecStart=/var/www/larynx/backend/venv/bin/python cleanup.py

[Install]
WantedBy=multi-user.target
